# scripts/picdown.service
[Unit]
Description=执行网络巡检操作，深信服，zabbix流量图抓取

[Service]
WorkingDirectory=/opt/picdown
ExecStart=/opt/picdown/bin/picdown